# mediaserver.service
[Unit]
Description=MediaServer
After=network.target

[Service]
User=brett
Group=brett
WorkingDirectory=/home/brett/Git/mediaserver
ExecStart=/home/brett/PyVEnvs/Py311/bin/python run.py mediaserver_config.yaml
Restart=always
RestartSec=30s
Type=simple

[Install]
WantedBy=multi-user.target